# voice/luna-voice.service
[Unit]
Description=Luna Voice Assistant
Wants=network-online.target
After=network-online.target

[Service]
User=aachten
WorkingDirectory=/home/aachten/luna-voice-assistant/voice
EnvironmentFile=-/home/aachten/luna-voice-assistant/voice/.env
ExecStart=/home/aachten/luna-voice-assistant/voice/venv/bin/python /home/aachten/luna-voice-assistant/voice/main.py
Restart=always
RestartSec=5s
Environment=PATH=/home/aachten/luna-voice-assistant/voice/venv/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
